# crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
SHELL=/bin/bash
* * * * * root cd /var/www/latest && php artisan schedule:run >> /dev/null 2>&1
